# doc/etc/profile.d/alias.sh
#!/bin/sh
# some more aliases
alias ll='ls -lh'
alias la='ls -A'
alias l='ls -CF'
alias calendar='calendar -A 3'
alias grep='grep --color'
alias indent='indent -kr -i8'
alias opera='opera -nomail'
