# .backup/etc/profile.d/jre.csh
#!/bin/csh
setenv JAVA_HOME /usr/lib/java
setenv MANPATH ${MANPATH}:/usr/lib/java/man
setenv PATH ${PATH}:/usr/lib/java/bin:/usr/lib/java/jre/bin
